# supervisor/dokuwiki.conf
[supervisord]
nodaemon=true

[program:phpfpm]
command= /etc/init.d/php7.4-fpm restart
autorestart= true

[program:nginx]
command= nginx -g "daemon off;" 
autorestart= true